# snapshot-backup.timer
[Unit]
Description=Daily SafetyTwin VM Snapshot Backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
